# rsync-watcher.path
[Unit]
Description=Monitorowanie katalogu upload

[Path]
PathChanged=%h/upload
PathModified=%h/upload
#PathExistsGlob=%h/upload/*

[Install]
WantedBy=default.target
