# rsync-watcher.service
[Unit]
Description=Rsync Runner Service

[Service]
Type=oneshot
ExecStart=%h/scripts/rsync-runner-pf.sh %f
