# rsync-watcher.path
[Unit]
Description=Directory monitoring

[Path]
PathChanged=%h/upload
PathModified=%h/upload

[Install]
WantedBy=default.target
